# crontab
# Facebook Daily Activity Cron Job
# Runs daily at 9 AM UTC (adjust timezone as needed)
0 9 * * * curl -X GET https://ul-cursor.onrender.com/api/cron/facebook-daily-activity
